# immich-dlna/rootfs/etc/services.d/immich-dlna/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

export IMMICH_URL
IMMICH_URL="$(bashio::config 'immich_url')"

export IMMICH_API_TOKEN
IMMICH_API_TOKEN="$(bashio::config 'immich_api_token')"

export IMMICH_VERIFY_SSL
IMMICH_VERIFY_SSL="${IMMICH_VERIFY_SSL:-true}"

export IMMICH_DLNA_HTTP_HOST
IMMICH_DLNA_HTTP_HOST="${IMMICH_DLNA_HTTP_HOST:-0.0.0.0}"

export IMMICH_DLNA_HTTP_PORT
IMMICH_DLNA_HTTP_PORT="$(bashio::config 'http_port')"

exec immich-dlna
