# immich-dlna/rootfs/etc/services.d/immich-dlna/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

export IMMICH_URL
IMMICH_URL="$(bashio::config 'immich_url')"

export IMMICH_API_TOKEN
IMMICH_API_TOKEN="$(bashio::config 'immich_api_token')"

export IMMICH_VERIFY_SSL
IMMICH_VERIFY_SSL="${IMMICH_VERIFY_SSL:-true}"

export IMMICH_DLNA_HTTP_HOST
IMMICH_DLNA_HTTP_HOST="${IMMICH_DLNA_HTTP_HOST:-0.0.0.0}"

export IMMICH_DLNA_HTTP_PORT
IMMICH_DLNA_HTTP_PORT="$(bashio::config 'http_port')"

# Ensure externally reachable media URLs for DLNA clients.
# If not explicitly provided, derive base URL from the first non-loopback IPv4.
if [ -z "${IMMICH_DLNA_BASE_URL:-}" ]; then
  detected_ip="$(hostname -I 2>/dev/null | awk '{for(i=1;i<=NF;i++) if ($i !~ /^127\./) {print $i; exit}}')"
  if [ -n "${detected_ip}" ]; then
    export IMMICH_DLNA_BASE_URL
    IMMICH_DLNA_BASE_URL="http://${detected_ip}:${IMMICH_DLNA_HTTP_PORT}"
  fi
fi

exec immich-dlna
